# src/systemd/smartlock.service

[Unit]
Description=OpenDDS SmartLock
After=network.target

[Service]
Type=simple
User=root
EnvironmentFile=/home/pi/smartlock/smartlock_env
ExecStart=/home/pi/smartlock/smartlock.sh start-system

[Install]
WantedBy=multi-user.target
